# archiso-hyprland/hyprland-zen/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f liveuser' --noclear --autologin liveuser %I $TERM
