# archiso-hyprland/hyprland-zen/airootfs/etc/skel/.config/systemd/user/wallpaper-changer.timer
[Unit]
Description=Run wallpaper changer every hour

[Timer]
OnBootSec=10s
OnUnitActiveSec=1h
RandomizedDelaySec=10s

[Install]
WantedBy=timers.target
